# systemd/system/lilith-ai.service
[Unit]
Description=Lilim AI Assistant Runtime (Rust Proxy & Python Brain)
Documentation=file:///usr/share/doc/lilith-ai/USER_MANUAL.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aegon
Group=aegon
ExecStart=/usr/bin/lilim-runtime
WorkingDirectory=/home/aegon/.local/share/lilim
Restart=on-failure
RestartSec=10
TimeoutStartSec=90
TimeoutStopSec=30

# Environment Configuration
Environment=LILIM_INSTALL=/usr/lib/lilim
Environment=LILIM_VENV=/usr/lib/lilim/venv
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/aegon
# Enable inference timing logs in journald
Environment=RUST_LOG=lilim_runtime=info,lilim_inference=info

# Security Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=no
ReadWritePaths=/home/aegon /var/log/lilim /tmp
ReadOnlyPaths=/etc/lilith /usr/lib/lilim
InaccessiblePaths=/root

# Additional Security
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
LockPersonality=yes
# NOTE: MemoryDenyWriteExecute removed — Python + Rust model inference requires JIT
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Resource Limits
# CPUQuota removed — inference is fully CPU-bound; throttling devastates tok/s.
# Candle GEMV reads ~1.7GB of weights per token; needs full CPU bandwidth.
# MemoryMax kept to prevent OOM; model + runtime sits at ~2.5GB RSS.
MemoryMax=4G
MemoryHigh=3.5G
TasksMax=1024

[Install]
WantedBy=multi-user.target
